# dockerize/supervisor/conf.d/mongodb.conf
[program:mongod]
command=/usr/bin/mongod --dbpath /mongo --port 27017 $MONGO_OPTIONS
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
