# dockerize/supervisor/conf.d/mongodb.conf
[program:mongod]
command=/usr/bin/mongod --dbpath /mongo --port 27017
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/mongod.log
